# rclone/rootfs/etc/services.d/rclone/finish
#!/usr/bin/env bashio

bashio::log.warning "Halt add-on"
exec /run/s6/basedir/bin/halt