# rclone/rootfs/etc/services.d/rclone/run
#!/usr/bin/with-contenv bashio

RCLONE_CONFIG=$(bashio::config 'configuration_path')
RCLONE_LOG=$(dirname $(bashio::config 'configuration_path'))'/sync.log'
REMOTES=$(bashio::config 'remotes')
LOCAL_RETENTION=$(bashio::config 'local_retention_days')

bashio::log.info "Pruning local files..."
find /backup/ -mtime +"${LOCAL_RETENTION}" -type f -delete
bashio::log.info "Pruning local files finished"

for remote in ${REMOTES}; do
    bashio::log.info $(bashio::jq "${remote}")
    remote_name=$(bashio::jq "${remote}" ".name")
    remote_path=$(bashio::jq "${remote}" ".path")
    remote_retention=$(bashio::jq "${remote}" ".retention_days")
    remote_start_url=$(bashio::jq "${remote}" ".start_url")
    remote_end_url=$(bashio::jq "${remote}" ".end_url")

    if [[ "${remote_start_url}" != "null" ]] 
    then
        bashio::log.info "Pinging start_url for ${remote_name}: ${remote_start_url}..."
        curl -fsS "${remote_start_url}"
    fi

    bashio::log.info "Starting remote copy for ${remote_name}..."

    rclone --config "${RCLONE_CONFIG}" \
        sync -v \
        --max-age "${remote_retention}d" \
        --delete-during \
        --delete-excluded \
        --ignore-errors \
        --ignore-existing /backup/ "${remote_name}:${remote_path}" \
        --log-level INFO \
        --log-file=${RCLONE_LOG}
    RET=$?

    bashio::log.info "Remote copy for ${remote_name} finished"

    if [[ "${remote_end_url}" != "null" ]]
    then
        bashio::log.info "Pinging end_url ${remote_name}: ${remote_end_url}"
        curl -fsS "${remote_end_url}/${RET}" --data-binary "@${RCLONE_LOG}"
    fi
done
