# victoriametrics/rootfs/etc/services.d/victoriametrics/run
#!/usr/bin/with-contenv bashio

RETENTION_PERIOD=$(bashio::config 'retention')
ADDITIONAL_ARGS=$(bashio::config 'additional_arguments')
CONFIG_PATH=$(bashio::config 'dynamic_configuration_path')

PROMETHEUS_ARGS=""
if [ -f "${CONFIG_PATH}/prometheus.yaml" ]; then
    PROMETHEUS_ARGS="-promscrape.config ${CONFIG_PATH}/prometheus.yaml"
fi

if bashio::config.false 'use_agent_instead'
then
    bashio::log.info "Starting Victoria Metrics Server"
    /usr/local/bin/victoria-metrics-prod -storageDataPath /data/victoria-metrics-data -retentionPeriod ${RETENTION_PERIOD} ${PROMETHEUS_ARGS} ${ADDITIONAL_ARGS}
else
    bashio::log.info "Starting Victoria Metrics Agent"
    /usr/local/bin/vmagent-prod -remoteWrite.tmpDataPath /data/vmagent-data -remoteWrite.url "$(bashio::config 'agent_remote_write_url')" ${PROMETHEUS_ARGS} ${ADDITIONAL_ARGS}
fi